# crontab
# daily at 0800
0 8 * * * python3 src/scraper.py notionloader | tee /logs/scraper.log
0 8 * * * python3 src/scraper.py companiesloader | tee /logs/scraper.log