# crontab
* * * * * python3 src/scraper.py load-notion | tee /logs/scraper.log